# ansible/templates/medusa-backend.service.j2
[Unit]
Description=Medusa Backend (API + Admin)
Documentation=https://docs.medusajs.com
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ app_dir }}/apps/backend/.medusa/server
EnvironmentFile={{ app_dir }}/apps/backend/.env
Environment=NODE_ENV={{ node_env }}
ExecStart=/home/nginx/.local/share/pnpm/pnpm start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=medusa-backend
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
